# root/etc/services.d/subfinder/run

#!/usr/bin/with-contenv bash

# 第一次启动全量扫描
s6-setuidgid subfinder /media \
	-c /config/subfinder.json

# 启动subfinder查找字幕
s6-setuidgid abc \
	looking4video
